# cncnet.service
[Unit]
Description=CnCNet Game Server (ARM Optimized)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cncnet
Group=cncnet
WorkingDirectory=/opt/cncnet

# Environment
Environment="RUST_LOG=info"
Environment="RUST_BACKTRACE=1"
Environment="TOKIO_WORKER_THREADS=2"

# Jemalloc tuning (statically linked via tikv-jemallocator)
Environment="MALLOC_CONF=background_thread:true,dirty_decay_ms:5000,muzzy_decay_ms:10000,lg_extent_max_active_fit:3"

# Memory and CPU limits for t4g.medium (4GB RAM, 2 vCPUs)
MemoryMax=3.5G
MemoryHigh=3G
CPUQuota=190%
CPUAffinity=0 1
TasksMax=4096

# IO and Network priority
IOWeight=100
IPAccounting=true

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cncnet/logs

# Performance tuning
LimitNOFILE=1048576
LimitNPROC=32768
LimitMEMLOCK=infinity

# Restart policy
Restart=always
RestartSec=5
StartLimitBurst=3
StartLimitInterval=60

# ARM CPU/network tuning (run as root before dropping to cncnet user)
ExecStartPre=-/bin/sh -c 'echo performance | tee /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor'
ExecStartPre=-/bin/sh -c 'echo 3 | tee /sys/class/net/*/queues/rx-*/rps_cpus'
ExecStartPre=-/bin/sh -c 'echo 2 | tee /proc/irq/*/smp_affinity'

# Main process
ExecStart=/opt/cncnet/cncnet-server \
    --port 50001 \
    --portv2 50000 \
    --name "CnCNet ARM Server" \
    --maxclients 200 \
    --iplimit 8 \
    --iplimitv2 4 \
    --workers 2 \
    --metrics-port 9090 \
    --log-level info

ExecReload=/bin/kill -USR1 $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

LimitCORE=infinity

[Install]
WantedBy=multi-user.target
